# dashboard/services/templates/python.service.j2
[Unit]
Description={{ app_name }} - Python Application
Documentation=https://github.com/quantyra/infrastructure
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ runtime_user }}
Group={{ runtime_group }}
WorkingDirectory={{ app_dir }}

# Environment setup
Environment="PATH={{ venv_path }}/bin:/usr/local/bin:/usr/bin:/bin"
Environment="VIRTUAL_ENV={{ venv_path }}"
Environment="APP_ENV={{ app_env }}"
{% if env_vars %}
{% for key, value in env_vars.items() %}
Environment="{{ key }}={{ value }}"
{% endfor %}
{% endif %}

# Start command - Gunicorn with production settings
ExecStart={{ venv_path }}/bin/gunicorn \
    --bind 0.0.0.0:{{ port }} \
    --workers {{ workers | default(4) }} \
    --threads {{ threads | default(2) }} \
    --worker-class {{ worker_class | default('sync') }} \
    --timeout {{ timeout | default(30) }} \
    --keep-alive {{ keep_alive | default(5) }} \
    --access-logfile {{ access_log | default('-') }} \
    --error-logfile {{ error_log | default('-') }} \
    --log-level {{ log_level | default('info') }} \
    {{ wsgi_app | default('app:app') }}

# Restart policy
Restart=always
RestartSec=5
StartLimitInterval=60
StartLimitBurst=3

# Resource limits
LimitNOFILE=65536
MemoryMax={{ memory_limit | default('512M') }}

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths={{ app_dir }}
{% if log_dir %}
ReadWritePaths={{ log_dir }}
{% endif %}

# Standard output/error handling
StandardOutput=journal
StandardError=journal
SyslogIdentifier={{ app_name }}

# Timeout settings
TimeoutStartSec=60
TimeoutStopSec=30

# Pre-start hook for Django migrations (optional)
# ExecStartPre={{ venv_path }}/bin/python manage.py migrate --noinput

[Install]
WantedBy=multi-user.target